# scripts/telecombot-backend-bot.service
[Unit]
Description=Telecombot backend, Telegram bot, and background workers
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/telecombot
EnvironmentFile=-/root/telecombot/.env
Environment=PYTHONUNBUFFERED=1
ExecStart=/root/telecombot/venv/bin/python /root/telecombot/run.py
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
